# mack/com.klp.honey-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.klp.honey-keepalive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/matthewjackson/Documents/KLP/honey/honey_keepalive.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>432000</integer><!-- every 5 days; free tier pauses after 7 days idle -->
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/honey_keepalive.err</string>
</dict>
</plist>
